# knot-resolver/knot-resolver-override.service
[Service]
ExecStart=
ExecStart=ip vrf exec vrf_external /usr/sbin/kresd --config=/etc/knot-resolver/kresd.conf
Sockets=
User=