# duplicity/files/ffmuc-backup.service
[Unit]
Description=FFMuc Backup

[Service]
ExecStart=/bin/bash /usr/local/sbin/backup.sh
