# network/ifupdown2/ff_fix_default_route.cron
#
# Check if a default route within vrf_external is configured and if it's active
#
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO={{ salt['pillar.get']('globals:ops_mail') }}

* * * * *       root    /usr/local/sbin/ff_fix_default_route
