# mongodb/mongodb_backup.cron
0 23  * * * root /usr/local/sbin/mongodb_backup
