# graylog/graylog-system-notifications.cron
0 12,20 * * * root /usr/local/sbin/graylog-system-notifications
