# graylog/graylog-group-mapping.cron
*/5 * * * * root /usr/local/sbin/graylog-group-mapping
