# apt/update_apt.cron
3 */6 * * * root apt-get update > /dev/null 2>&1
