# batman/ff_check_gateway.cron
#
# Check if THE INTERNET is reachable via chosen upstream (Salt managed)
#
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO={{ salt['pillar.get']("ffho:tech_c_mail") }}

* * * * *       root    /usr/local/sbin/ff_check_gateway
