# pppoe/ff_fix_ppp.cron
*/5 *  * * *	root	/usr/local/sbin/ff_fix_ppp
